# daspi-services/daspi-apps-usage.service
[Unit]
Description=DasPi initial server
After=network.target

[Service]
ExecStart=/usr/bin/python3 /home/ubuntu/DarPi/backend/save_apps.py

RestartSec=5
Restart=always

StandardError=file:/home/ubuntu/logs/daspi-apps-stderr.log

[Install]
WantedBy=multi-user.target
